# roles/basic/files/init/boot-log.conf
description "Log shutdown event"

start on started mongodb

task
chdir /opt/eventlog

script
	./events.py add Information "Boot"
end script 
