# roles/basic/files/init/dnstest.conf
description "AppChord DNS Test"

start on runlevel [234]
stop on runlevel [!234]

chdir /usr/local/bin
exec python dnstest.py www.google.com 80
respawn
