# roles/basic/files/init/shutdown-log.conf
description "Log shutdown event"

start on stopping mongodb

task
chdir /opt/eventlog

script
	./events.py add Information "Shutdown"
end script 
