# roles/kiosk/files/kiosk.conf
 # kiosk.conf

description "AppChord Kiosk App"

start on started networking
stop on runlevel [016]


limit nofile 1000000 1000000


script

  cd /var/www/kiosk
  npm start

end script

respawn
